# scripts/com.weft.session-archive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.weft.session-archive</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/rhhart</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bun</string>
        <string>/Users/rhhart/Documents/GitHub/weft-dev/scripts/session-archive.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rhhart/Documents/GitHub/weft-dev</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/rhhart/.config/weft/session-archive.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rhhart/.config/weft/session-archive.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
